# systemd/aceproxy.service
[Unit]
Description=AceProxy
After=network.target

[Service]
User=nobody
ExecStart=/opt/aceproxy/acehttp.py
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target 
